# ods-parenting-exchange/ods-generic-parent.service
[Unit]
Description=DNSSEC Generic Parenting
After=rabbitmq-server.service

[Install]

[Service]
Type = simple
#TODO#HOME#
Environment = PYTHONPATH=@PREFIX@/lib
ExecStart = @PREFIX@/bin/ods-generic-parent

